# supervisord.conf
[supervisord]
nodaemon=true
 
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
 
[program:jenkins]
command=service jenkins start
autostart=true
autorestart=true